# puppet/environments/shared/modules/domotd/files/rc.local
#!/bin/bash
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

# touch the file if it exists
[ -f /var/lock/subsys/local ] && touch /var/lock/subsys/local
